# deploy/agenticmojo.simple.service
[Unit]
Description=AgenticMojo (FastAPI/uvicorn)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/agenticmojo
Environment=HOST=127.0.0.1
Environment=PORT=8000
# Optional: point to your config file if customized
# Environment=TRAE_CONFIG_FILE=/opt/agenticmojo/trae_config.yaml
ExecStart=/opt/agenticmojo/.venv/bin/python -m uvicorn server.app:app --host 127.0.0.1 --port 8000 --no-access-log
Restart=always
RestartSec=2
; User and Group will be set by the deploy script when installing
; User=
; Group=

[Install]
WantedBy=multi-user.target
